# dvac.crontab

# backup script at 00:00 AM everyday
0 0 * * * /var/www/dvac/backup/backup.sh

# proceedings checker at 10:00 AM everyday
0 10 * * *	php /var/www/dvac/sms/proceedings_lookahead.php > /tmp/proceedings.log
