# backup/backup.crontab
0 0 * * * /home/krupa/work/dvac/backup/backup.sh
